# raspberry/scripts/systemd/savonia-iot-sensor@.service
[Unit]
Description=Savonia IoT - Sensor Reader (%i)
Wants=time-sync.target
After=time-sync.target

# If systemd-time-wait-sync.service exists on the system, waiting for it is the most reliable
# way to ensure the clock is actually synchronized before reading sensors.
After=systemd-time-wait-sync.service
Wants=systemd-time-wait-sync.service


[Service]
EnvironmentFile=/etc/savonia-iot/env
Type=oneshot

# Run as the installing user (set by installer); fallback is fine for manual runs
User=%i

# Set a stable working directory (installer ensures ownership)
WorkingDirectory=/opt/savonia-iot

# Use env to resolve node from PATH (works for /usr/bin/node and /usr/local/bin/node)
ExecStart=/usr/bin/env node raspberry/dist/sensor-reader.js --config raspberry/config/config.json --sensorId %i

# Ensure PATH includes common Node locations
Environment=PATH=/usr/local/bin:/usr/bin:/bin
